# deploy/systemd/wethr-calibration-monthly.timer
[Unit]
Description=Run Wethr archive and candidate training monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target
