# deploy/systemd/wethr-calibration-daily.timer
[Unit]
Description=Run Wethr calibration collection daily

[Timer]
OnCalendar=*-*-* 04:15:00 UTC
Persistent=true

[Install]
WantedBy=timers.target
