# deploy/systemd/wethr-export.service
[Unit]
Description=Export Wethr settled trades for n8n audit

[Service]
Type=oneshot
WorkingDirectory=%h/projects/wethr/collector
ExecStart=/usr/bin/python3 %h/projects/wethr/collector/run.py export-settled
